# debian/vyatta-vci-bus.vyatta-vci-bus.socket
[Unit]
Description=D-Bus VCI Message Bus Socket

[Socket]
ListenStream=/var/run/vci/vci_bus_socket

[Install]
WantedBy=sockets.target
